# scripts/systemd/twitter-memory.service.template
[Unit]
Description=twitter-memory local ingest backend
After=network.target

[Service]
Type=simple
WorkingDirectory={{REPO}}
EnvironmentFile={{ENV_FILE}}
ExecStart={{VENV_PYTHON}} -m backend.main
Restart=on-failure
RestartSec=3
# Logs go to journalctl --user -u twitter-memory -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
